# etc/supervisor/conf.d/cron.conf
[program:cron]
directory = /etc/cron.d/
command = /usr/sbin/cron -f
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autorestart = true
