# etc/supervisor/conf.d/denyhosts.conf
[program:denyhosts]
directory = /usr/sbin/
command = /usr/sbin/denyhosts --foreground
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autorestart = true
