# etc/supervisor/conf.d/memcached.conf
[program:memcached]
autorestart=true
command=/usr/bin/memcached -u memcache -l 127.0.0.1 -m 8 -p 11211 -c 1024 -vv
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
