# etc/supervisor/conf.d/php-fpm56.conf
[program:php-fpm5.6]
command=/usr/sbin/php-fpm5.6 --nodaemonize
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
