# etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server --daemonize no
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
