# etc/supervisor/conf.d/tor.conf
[program:tor]
autorestart = true
command=/usr/sbin/tor
redirect_stderr = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
